# oakd-recorder.service
[Unit]
Description=OAK-D Hand Tracking Recorder
After=network.target

[Service]
Type=simple
User=asimov
WorkingDirectory=/home/asimov/pi-recorder
ExecStart=/home/asimov/pi-recorder/venv/bin/python /home/asimov/pi-recorder/record.py --recordings-dir /home/asimov/pi-recorder/recordings
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Give it time to initialize OAK-D
TimeoutStartSec=60

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
